# ansible/roles/backup/templates/backup.timer.j2
[Unit]
Description=Runs {{ game_name }} backup service periodically

[Timer]
OnCalendar=hourly
persistent=true

[Install]
WantedBy=timers.target
